# ~/.local/etc/cron.daily/sync-codes.sh
#!/bin/bash

rsync -av --exclude='node_modules' ~/Dropbox/Shared/Computer_code/ ~/Documents/Computer_code
